# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Default queue: general app jobs (emails, reports, etc.)
[program:queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work database --queue=default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
numprocs=1
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/queue-default.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/queue-default-error.log

; FBR fiscal outbox queue: dedicated workers so a slow/down FBR endpoint
; never starves other app jobs. No --tries override here: FiscalizeInvoiceJob
; sets its own $tries=100 safety ceiling and drives retry timing itself via
; $this->release($delay) computed from fiscal_outbox's exponential backoff -
; the real give-up threshold is fiscal.max_retry_attempts, enforced by
; FiscalSubmissionService, not by the queue worker's attempt counter.
[program:queue-fiscal]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work database --queue=fiscal --sleep=2 --max-time=3600 --timeout=90
autostart=true
autorestart=true
numprocs=2
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/queue-fiscal.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/queue-fiscal-error.log

[program:scheduler]
command=/bin/sh -c "while true; do php /var/www/html/artisan schedule:run --verbose --no-interaction; sleep 60; done"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/scheduler.log
stderr_logfile=/var/log/supervisor/scheduler-error.log
